# ndm-redis1/crontab/crontab.conf
* * * * * /etc/keepalived/scripts/check_keepalived.sh
